# supervisor-schedmon.conf
[program:schedmon]
user=mupif
command=ttyd --readonly --port 8001 /usr/bin/python3 -m mupif.cli schedmon
