# supervisor-proxy.conf
[program:apache2]
directory=%(ENV_MUPIF_HOME_DIR)s
environment=APACHE_RUN_DIR=%(ENV_MUPIF_HOME_DIR)s
user=root
command=/usr/sbin/apache2 -d. -f ./apache2.conf -D FOREGROUND
